# resource/supervisord.conf
[supervisord]
nodaemon=true
[program:consul]
command=/bin/bash -c "consul agent -join $CONSUL_HOST -data-dir /data/consul -config-dir /etc/consul.d"

